# services/inotify/run
#!/bin/sh

while /sbin/setuser abc inotifywait -e modify /config/cron/crontab; do
  cp /config/cron/crontab /etc/cron.d/
  echo "New cron file detected, copying"
  restart cron

done