# root/etc/services.d/plexemail/run
#!/usr/bin/with-contenv bash

cd /app/plexemail/web/

exec \
	s6-setuidgid abc python -m SimpleHTTPServer 8080
