# services/python_web/run
#!/bin/bash
cd /app/plexemail/web/
/sbin/setuser abc python -m SimpleHTTPServer 80
